# systemd/cloudflared-canvas@.service
[Unit]
Description=Cloudflare Tunnel for TG Canvas (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=joker
Group=joker
Environment=HOME=/home/joker
Environment=TG_CANVAS_INSTANCE=%i
EnvironmentFile=/etc/tg-canvas/%i.env
ExecStartPre=/usr/bin/mkdir -p /home/joker/projects/tg-canvas-studio/logs
ExecStart=/bin/bash -lc '/usr/local/bin/cloudflared tunnel run ${CLOUDFLARED_TUNNEL:?set_CLOUDFLARED_TUNNEL_in_/etc/tg-canvas/%i.env}'
Restart=always
RestartSec=5s
StandardOutput=append:/home/joker/projects/tg-canvas-studio/logs/cloudflared-canvas-%i.log
StandardError=append:/home/joker/projects/tg-canvas-studio/logs/cloudflared-canvas-%i.log

[Install]
WantedBy=multi-user.target
